# debian/cblues-lxqt-livesettings.service
[Unit]
Description=LXQt cheatcodes in live media
ConditionPathExists=/usr/share/cblues-lxqt-livesettings/
After=cblues-common-livesettings.service

[Service]
Type=oneshot
ExecStart=/usr/share/cblues-lxqt-livesettings/lxqt-cheatcodes start

[Install]
WantedBy=graphical.target
