# debian/xperimental-lxqt-livesettings.service
[Unit]
Description=LXQt cheatcodes in live media
ConditionPathExists=/usr/share/xperimental-lxqt-livesettings/
After=xperimental-common-livesettings.service

[Service]
Type=oneshot
ExecStart=/usr/share/xperimental-lxqt-livesettings/lxqt-cheatcodes start

[Install]
WantedBy=graphical.target
